# upstart/ronin.conf
description "ronin os"
author "Phoxygen"
start on runlevel [2345]
exec startx /opt/b2g/session.sh
